# packages/gbp-webhook/gbp_webhook-2.0.0.tar.gz/gbp_webhook-2.0.0/src/gbp_webhook/templates/gbp-webhook.service
[Unit]
Description=gbp-notifications webhook responder
Wants=dbus.socket
After=dbus.socket
After=graphical-session.target

[Service]
EnvironmentFile={{ config_path }}
ExecStart={{ gbp_path }} webhook serve $GBP_WEBHOOK_ARGS

[Install]
WantedBy=default.target
